# systemd/daystrom.service
[Unit]
Description=Daystrom AI Assistant
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=daystrom
Group=daystrom
WorkingDirectory=/home/daystrom/app
Environment="PATH=/home/daystrom/app/venv/bin"
ExecStart=/home/daystrom/app/venv/bin/python main.py
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/daystrom/app/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daystrom

[Install]
WantedBy=multi-user.target

